# systemd/mwu-bot-aur.service
[Unit]
Description=MWU trading bot scheduler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/mwu-bot
ExecStart=/usr/bin/python -m bot.scheduler
Restart=on-failure
RestartSec=10
User=mwu
Group=mwu
UMask=007

[Install]
WantedBy=multi-user.target

